# dot/profile.d/coding.sh
export CODING_HOME="$HOME/Coding"
export CDPATH=$CDPATH:$CODING_HOME
