# dot/profile.d/gaming.sh
export GAME_HOME=$HOME/Games
export CDPATH=$CDPATH:$GAME_HOME
