# dot/profile.d/lua.sh
test -r $HOME/.luavmrc && . $HOME/.luavmrc

alias rock='luarocks'
alias ilua='lua -l rlcompleter'

