# dot/profile.d/ruby.sh
export PATH=$PATH:$(ruby -rubygems -e "puts Gem.user_dir")/bin
